# .config/systemd/user/gpg-agent.service.d/override.conf
[Service]
Environment="QT_QPA_PLATFORMTHEME=qt6ct"
Environment="QT_STYLE_OVERRIDE=Fusion"
